# cooling-hat-oled/cooling.service
[Unit]
Description=Cooling service
After=network.target

[Service]
ExecStart=/opt/cooling-fan/startup.sh

[Install]
WantedBy=multi-user.target